# redeclipse_server.cron
*/1 * * * * /etc/rc.d/rc.redeclipse_server check
